# package/etc/systemd/system/maulu.service
[Unit]
Description=A simple URL shortening backend.

[Service]
ExecStart=/usr/bin/maulu

[Install]
WantedBy=multi-user.target
